# links/.config/systemd/user/prlcc.service
[Unit]
Description=Parallels Control Center

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/prlcc
