# links/.config/systemd/user/prlcp.service
[Unit]
Description=Parallels Copy-Paste Tool

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/prlcp
